# cron/crontab
MAILTO=""
@reboot     loxberry cd REPLACELBPBINDIR && REPLACELBPBINDIR/watchdog.sh > /dev/null 2>&1
